# artifacts/08-infra/crontab
# ClawBTP — crontab supercronic
# Tous les endpoints sont sécurisés par le header x-cron-secret
# APP_INTERNAL_URL et CRON_SECRET sont injectés via Docker secrets au runtime
#
# Ordre UTC chaque lundi matin (D-7-03 BINDING) :
#   06h00 — jalons (vérification toutes les heures, inclut 06h00)
#   06h30 — briefing prospectif (Sprint 7)
#   07h00 — détection dérives (Sprint 6)
#   07h15 — rapports hebdomadaires (Sprint 5)
# Les crons quotidiens/horaires sont hors de cet ordre spécifique lundi.

# Génération CRs journaliers — tous les jours à 18h00 UTC
0 18 * * * curl -sf -X POST -H "x-cron-secret: ${CRON_SECRET}" "${APP_INTERNAL_URL}/api/cron/cr"

# Vérification jalons — toutes les heures (inclut lundi 06h00 UTC)
0 * * * * curl -sf -X POST -H "x-cron-secret: ${CRON_SECRET}" "${APP_INTERNAL_URL}/api/cron/jalons"

# Briefing prospectif lundi matin — lundi 06h30 UTC (D-7-03 / 30 6 * * 1)
# AVANT dérives (07h00) et hebdo (07h15) — ordre binding architecture Sprint 7
30 6 * * 1 curl -sf -X POST -H "x-cron-secret: ${CRON_SECRET}" "${APP_INTERNAL_URL}/api/cron/briefing"

# Détection proactive des dérives — tous les jours à 07h00 UTC (D-6-02 / ADR-6-005 / D-6-13)
# Endpoint canonique : /api/cron/derives (PLURIEL) — corrige la ligne singulière morte héritée Sprint 4
0 7 * * * curl -sf -X POST -H "x-cron-secret: ${CRON_SECRET}" "${APP_INTERNAL_URL}/api/cron/derives"

# Génération rapports hebdomadaires — tous les lundis à 07h15 UTC (AM-01 BINDING : ligne DISTINCTE)
15 7 * * 1 curl -sf -X POST -H "x-cron-secret: ${CRON_SECRET}" "${APP_INTERNAL_URL}/api/cron/rapports-hebdo"
